# files/etc/uci-defaults/99_install
#!/bin/sh

opkg install /root/luci-i18n-adguardhome-zh-cn_1.12-1_all.ipk

rm -f /root/luci-i18n-adguardhome-zh-cn_1.12-1_all.ipk

exit 0
